# pfPorts/isc-dhcp42-server/pkg-plist
@comment $FreeBSD: head/net/isc-dhcp42-server/pkg-plist 316014 2013-04-18 01:25:41Z zi $
@unexec if cmp -s %D/etc/dhcpd.conf.example %D/etc/dhcpd.conf; then rm -f %D/etc/dhcpd.conf; fi
@unexec echo "You should stop the dhcpd and dhcpd6 daemons if you are permanently removing this port."
etc/dhcpd.conf.example
@exec if [ ! -f %D/etc/dhcpd.conf ] ; then cp -p %D/%F %B/dhcpd.conf; fi
%%IPV6%%@unexec if cmp -s %D/etc/dhcpd6.conf.example %D/etc/dhcpd6.conf; then rm -f %D/etc/dhcpd6.conf; fi
%%IPV6%%etc/dhcpd6.conf.example
%%IPV6%%@exec if [ ! -f %D/etc/dhcpd6.conf ] ; then cp -p %D/%F %B/dhcpd6.conf; fi
%%IPV6%%etc/rc.d/isc-dhcpd6
sbin/dhcpd
bin/omshell
%%LDAP%%bin/dhcpd-conf-to-ldap
include/dhcpctl.h
include/isc-dhcp/dst.h
include/isc-dhcp/app.h
include/isc-dhcp/assertions.h
include/isc-dhcp/atomic.h
include/isc-dhcp/base64.h
include/isc-dhcp/bind9.h
include/isc-dhcp/bitstring.h
include/isc-dhcp/boolean.h
include/isc-dhcp/buffer.h
include/isc-dhcp/bufferlist.h
include/isc-dhcp/commandline.h
include/isc-dhcp/condition.h
include/isc-dhcp/dir.h
include/isc-dhcp/entropy.h
include/isc-dhcp/error.h
include/isc-dhcp/event.h
include/isc-dhcp/eventclass.h
include/isc-dhcp/file.h
include/isc-dhcp/formatcheck.h
include/isc-dhcp/fsaccess.h
include/isc-dhcp/hash.h
include/isc-dhcp/heap.h
include/isc-dhcp/hex.h
include/isc-dhcp/hmacmd5.h
include/isc-dhcp/httpd.h
include/isc-dhcp/int.h
include/isc-dhcp/interfaceiter.h
include/isc-dhcp/iterated_hash.h
include/isc-dhcp/lang.h
include/isc-dhcp/lex.h
include/isc-dhcp/lfsr.h
include/isc-dhcp/lib.h
include/isc-dhcp/list.h
include/isc-dhcp/log.h
include/isc-dhcp/magic.h
include/isc-dhcp/md5.h
include/isc-dhcp/mem.h
include/isc-dhcp/msgcat.h
include/isc-dhcp/msgs.h
include/isc-dhcp/mutex.h
include/isc-dhcp/mutexblock.h
include/isc-dhcp/namespace.h
include/isc-dhcp/net.h
include/isc-dhcp/netaddr.h
include/isc-dhcp/netdb.h
include/isc-dhcp/offset.h
include/isc-dhcp/once.h
include/isc-dhcp/ondestroy.h
include/isc-dhcp/os.h
include/isc-dhcp/parseint.h
include/isc-dhcp/platform.h
include/isc-dhcp/print.h
include/isc-dhcp/quota.h
include/isc-dhcp/radix.h
include/isc-dhcp/random.h
include/isc-dhcp/ratelimiter.h
include/isc-dhcp/refcount.h
include/isc-dhcp/region.h
include/isc-dhcp/resource.h
include/isc-dhcp/result.h
include/isc-dhcp/resultclass.h
include/isc-dhcp/rwlock.h
include/isc-dhcp/serial.h
include/isc-dhcp/sha1.h
include/isc-dhcp/sha2.h
include/isc-dhcp/sockaddr.h
include/isc-dhcp/socket.h
include/isc-dhcp/stdio.h
include/isc-dhcp/stdlib.h
include/isc-dhcp/stdtime.h
include/isc-dhcp/string.h
include/isc-dhcp/symtab.h
include/isc-dhcp/syslog.h
include/isc-dhcp/task.h
include/isc-dhcp/taskpool.h
include/isc-dhcp/thread.h
include/isc-dhcp/time.h
include/isc-dhcp/timer.h
include/isc-dhcp/types.h
include/isc-dhcp/util.h
include/isc-dhcp/version.h
include/isc-dhcp/xml.h
include/omapip/alloc.h
include/omapip/buffer.h
include/omapip/omapip.h
include/omapip/omapip_p.h
include/omapip/hash.h
include/omapip/trace.h
include/omapip/convert.h
include/omapip/isclib.h
include/omapip/result.h
lib/libdhcpctl.a
lib/libomapi.a
lib/libdst.a
lib/libdns-dhcp.a
lib/libisc-dhcp.a
%%LDAP%%%%SCHEMA_DIR%%/dhcp.schema
%%LDAP%%@dirrm %%SCHEMA_DIR%%
@dirrm include/omapip
@dirrm include/isc-dhcp
